# cubesat/rc.local
#!/bin/bash
python3 -u /home/raspberrypi/cubesat-2025/cubesat/control_loop.py &> /home/raspberrypi/cubesat-2025/cubesat/log.txt &
exit 0
